# scripts/com.user.claude-continue.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.claude-continue.__STATE_ID__</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>osascript -e 'tell application "__TERM_APP__" to do script "cd \"__PROJECT_PATH__\" &amp;&amp; CLAUDE_CONTINUE_UNATTENDED=1 claude -c \"/awake\""'</string>
  </array>
  <key>StartInterval</key>
  <integer>18000</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/claude-continue-__STATE_ID__.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude-continue-__STATE_ID__.err</string>
</dict>
</plist>
